# systemfiles/controller.service
[Unit]
Description=Airiana Controller TCP(http) --> UDP server
After=network-online.target
Type=idle
[Service]
ExecStart=/usr/bin/python3 /home/pi/airiana/public/controler.py
Restart=always
serviceUser=root
[Install]
WantedBy=multi-user.target
