# systemfiles/controller.service
[Unit]
Description=Airiana Controller TCP(http) --> UDP server
After=network-online.target
Type=idle
[Service]
ExecStart=/home/pi/airiana/venv/bin/python  /home/pi/airiana/public/controller.py
Restart=always
serviceUser=root
[Install]
WantedBy=multi-user.target
